# assets/configs/etc/systemd/system/avd-bridge.socket
[Socket]
ListenStream=/run/avd-bridge/avd.sock
SocketMode=0666
RuntimeDirectory=avd-bridge

[Install]
WantedBy=sockets.target
