# assets/configs/etc/systemd/system/fast-charge.timer
[Unit]
Description=Enable fast charging at 07:00

[Timer]
OnCalendar=Mon..Fri *-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
