# assets/configs/etc/systemd/system/fast-charge.service
[Unit]
Description=Enable fast charging

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'echo Fast > /sys/class/power_supply/BAT0/charge_types'
