# assets/configs/etc/systemd/system/mnt-storagebox.automount
[Unit]
Description=Automount Hetzner Storage Box (devcanister/files)

[Automount]
Where=/mnt/storagebox
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
