# assets/configs/etc/systemd/system/ollama-bridge.socket
[Unit]
Description=Unix socket bridge to Ollama

[Socket]
ListenStream=/run/ollama-bridge/ollama.sock
SocketMode=0660
RuntimeDirectory=ollama-bridge

[Install]
WantedBy=sockets.target
