# assets/configs/etc/systemd/system/ollama-bridge-client.service
[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd /run/ollama-bridge/ollama.sock
