# assets/configs/etc/systemd/system/platform-profile-tlp.service
[Unit]
Description=Apply TLP profile from platform_profile

[Service]
Type=oneshot
ExecStart=/usr/local/bin/platform-profile-to-tlp
