# assets/configs/etc/systemd/system/forgejo-runner.service.d/override.conf
[Service]
Environment="DOCKER_HOST=unix:///run/user/941/podman/podman.sock"
CPUQuota=200%
