# assets/configs/etc/systemd/system/forgejo-runner.service.d/override.conf
[Service]
Environment="DOCKER_HOST=unix:///run/user/949/podman/podman.sock"
LoadCredentialEncrypted=forgejo-personal
CPUQuota=1000%
